# docker-tests/docker/service-supervisords/druid.conf
[program:druid-service]
command=java %(ENV_DRUID_SERVICE_JAVA_OPTS)s %(ENV_DRUID_COMMON_JAVA_OPTS)s -cp %(ENV_DRUID_COMMON_CONF_DIR)s:%(ENV_DRUID_SERVICE_CONF_DIR)s:%(ENV_DRUID_DEP_LIB_DIR)s org.apache.druid.cli.Main server %(ENV_DRUID_SERVICE)s
redirect_stderr=true
priority=100
autorestart=false
stdout_logfile=/shared/logs/%(ENV_DRUID_LOG_NAME)s
